# resmon.service
[Unit]
Description=ResMon

[Service]
ExecStart=/opt/resmon/startup/resmon.sh
Restart=always
